# recipes-support/swupdate/swupdate/swupdate-check.timer
[Unit]
Description=Run swupdate-check.service every 12 hours

[Timer]
OnBootSec=15min
OnUnitActiveSec=12h

[Install]
WantedBy=timers.target
