# recipes-support/swupdate/swupdate/swupdate-check.service
[Unit]
Description=Fetch new Linux system images, install if applicable"

[Service]
ExecStart=/usr/bin/swupdate-check
Type=oneshot
